# readme-sync.service
[Unit]
Description=README Sync Pipeline - Webhook Listener (Multi-Repo)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/git-readme-sync-sword
EnvironmentFile=/home/ubuntu/git-readme-sync-sword/.env
ExecStart=/home/ubuntu/git-readme-sync-sword/venv/bin/uvicorn \
    app:app \
    --host 0.0.0.0 \
    --port 5000 \
    --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target